# salt/schedule-commands/schedule-commands.service
[Unit]
Description=Salt minion command scheduler
After=network-online.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/opt/ops-tools/salt/schedule-commands.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
