# backups/psql-backup/psql-backup.cron
35 16 * * * root /opt/opt-tools/backups/psql-backup/psql-backup.py